# build/systemd/getty@tty2.service
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin monitoring --noissue --noclear %I $TERM